# clobberer/cleanup_clobberer.cron
MAILTO=catlee@mozilla.com
@hourly CLEANUP_PASSWORD=sekrits /usr/bin/php /var/www/html/build/clobberer/cleanup.php
